# root/etc/s6-overlay/s6-rc.d/init-rdesktop/run
#!/usr/bin/with-contenv bash

mkdir -p /var/run/xrdp || exit 1
chown root:xrdp /var/run/xrdp || exit 1
chmod 2775 /var/run/xrdp || exit 1

mkdir -p /var/run/xrdp/sockdir || exit 1
chown root:xrdp /var/run/xrdp/sockdir || exit 1
chmod 3777 /var/run/xrdp/sockdir || exit 1

if [ ! -f "/keylock" ]; then
    cd /etc/xrdp
    xrdp-keygen xrdp
    rm -f /etc/xrdp/*.pem
    openssl req -x509 -nodes -days 3650 -newkey rsa:2048 \
    -keyout /etc/xrdp/key.pem \
    -out /etc/xrdp/cert.pem \
    -subj "/C=US/ST=CA/L=Carlsbad/O=Linuxserver.io/OU=LSIO Server/CN=*"
    touch /keylock
fi

if [ ! -e /lock.file ]; then
  # give abc a sudo shell
  chsh abc -s /bin/bash
  sed -e 's/%sudo	ALL=(ALL:ALL) ALL/%sudo ALL=(ALL:ALL) NOPASSWD: ALL/g' \
    -i /etc/sudoers
  sed -e 's/^wheel:\(.*\)/wheel:\1,abc/g' -i /etc/group
fi

# create lock file after first run
touch /lock.file

# default file copies first run
[[ ! -d /config/.config ]] && \
	mkdir -p /config/.config && \
        cp /defaults/bashrc /config/.bashrc && \
	cp /defaults/startwm.sh /config/startwm.sh

# permissions
PERM=$(stat -c '%U' /config/.config)
[[ "${PERM}" != "abc" ]] && \
        chown -R abc:abc /config
